# services/dbus-mqtt-chain2/run
#!/bin/sh
exec 2>&1

MQTT_BROKER="127.0.0.1"
BATTERIES="4"

cd /data/dbus-mqtt-battery
exec python3 dbus-mqtt-battery.py \
    --broker "$MQTT_BROKER" \
    --batteries "$BATTERIES" \
    --instance 513 \
    --topic-prefix battery2 \
    --service-suffix mqtt_chain2 \
    --product-name "JBD Battery Chain 2"
